# gitmirror.timer
[Unit]
Description=Daily locate database update

[Timer]
OnCalendar=*:0/5
Persistent=true
